# systemd/gb10-stack-recovery.service
[Unit]
Description=GB10 Tier 2 coordinated model stack recovery
After=gb10-memory-guardian.service
StartLimitIntervalSec=0

[Service]
Type=oneshot
Environment=GB10_STACK_RECOVERY_MIN_MEM_AVAILABLE_GIB=40
Environment=GB10_STACK_RECOVERY_COMMAND_TIMEOUT_SECONDS=15
Environment=GB10_STACK_RECOVERY_STOP_TIMEOUT_SECONDS=90
Environment=GB10_STACK_RECOVERY_START_TIMEOUT_SECONDS=900
Environment=GB10_STACK_RECOVERY_READINESS_TIMEOUT_SECONDS=900
Environment=GB10_STACK_RECOVERY_PROBE_TIMEOUT_SECONDS=5
Environment=GB10_STACK_RECOVERY_POLL_SECONDS=5
Environment=GB10_STACK_RECOVERY_KILL_AFTER_SECONDS=5
ExecStart=%h/.local/bin/gb10_stack_recovery.sh
TimeoutStartSec=7200
TimeoutStopSec=10
KillMode=control-group
MemoryMax=128M
MemorySwapMax=0
OOMScoreAdjust=-900
RuntimeDirectory=gb10-stack-recovery
RuntimeDirectoryMode=0700
StateDirectory=gb10-stack-recovery
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
StandardOutput=journal
StandardError=journal
