# systemd/gb10-memory-guardian.service
[Unit]
Description=GB10 pre-reserved direct cgroup memory guardian
After=default.target
StartLimitIntervalSec=0

[Service]
Type=simple
# [thresholds] in config.toml overrides these environment defaults at runtime.
Environment=GB10_MEMORY_GUARDIAN_RESERVE_MIB=64
Environment=GB10_MEMORY_GUARDIAN_MEM_AVAIL_STOP_GIB=1
Environment=GB10_MEMORY_GUARDIAN_POLL_SECONDS=1
Environment=GB10_MEMORY_GUARDIAN_RETRY_SECONDS=5
Environment=GB10_MEMORY_GUARDIAN_ESCALATION_GRACE_SECONDS=60
Environment=GB10_MEMORY_GUARDIAN_ESCALATION_UNIT=gb10-stack-recovery.service
Environment=GB10_MEMORY_GUARDIAN_ESCALATION_ENABLED=false
Environment=GB10_MEMORY_GUARDIAN_CONFIG_PATH=%h/.config/gb10-memory-guardian/config.toml
Environment=GB10_MEMORY_GUARDIAN_EXPECTED_LABEL=aeon-text
Environment=GB10_MEMORY_GUARDIAN_EXPECTED_REGISTRATION_FILE=text-cgroup.v1
ExecStart=%h/.local/bin/gb10-memory-guardian
Restart=always
RestartSec=1
TimeoutStopSec=5
MemoryMin=64M
MemoryMax=96M
MemorySwapMax=0
OOMScoreAdjust=100
RuntimeDirectory=gb10-memory-guardian
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes

[Install]
WantedBy=default.target
