# systemd/gb10-swap-guard.service
[Unit]
Description=GB10 vLLM swap/free-memory guard (1s detection; stop reranker on low MemAvailable or high swap)
After=default.target
StartLimitIntervalSec=0

[Service]
Type=simple
Environment=GB10_SWAP_WARN_GIB=7.5
Environment=GB10_SWAP_STOP_GIB=12
Environment=GB10_MEM_AVAIL_STOP_GIB=1
Environment=GB10_SWAP_GUARD_INTERVAL=1
Environment=GB10_SWAP_GUARD_SAMPLE_LOG_INTERVAL=20
Environment=GB10_SWAP_ATTRIBUTION_INTERVAL=300
ExecStart=%h/.local/bin/gb10-swap-guard.sh
Restart=always
RestartSec=5
MemoryMax=128M
MemorySwapMax=0
OOMScoreAdjust=-850

[Install]
WantedBy=default.target
