# systemd/gb10-memory-guardian.service
[Unit]
Description=GB10 pre-reserved direct cgroup memory guardian
After=default.target
StartLimitIntervalSec=0

[Service]
Type=simple
Environment=GB10_MEMORY_GUARDIAN_RESERVE_MIB=64
Environment=GB10_MEMORY_GUARDIAN_MEM_AVAIL_STOP_GIB=1
Environment=GB10_MEMORY_GUARDIAN_POLL_SECONDS=1
Environment=GB10_MEMORY_GUARDIAN_RETRY_SECONDS=5
ExecStart=%h/.local/bin/gb10-memory-guardian
Restart=always
RestartSec=1
TimeoutStopSec=5
MemoryMin=64M
MemoryMax=96M
MemorySwapMax=0
OOMScoreAdjust=-1000
RuntimeDirectory=gb10-memory-guardian
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectHome=read-only
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes

[Install]
WantedBy=default.target
