# systemd/sysmon.service
[Unit]
# Observer-only telemetry. Recovery remains exclusively owned by the Rust guardian.
Description=Observer-only system monitor (1s target interval; measured cadence)
After=default.target
# Prevent systemd from ever giving up on restarts.
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=%h/.local/bin/sysmon.sh
Restart=always
RestartSec=5
# Make OOM killer strongly prefer other processes over sysmon
OOMScoreAdjust=-900

[Install]
WantedBy=default.target
