# systemd/vllm-aeon-27b-dflash.service
[Unit]
Description=vLLM AEON-7 Qwen3.6-27B DFlash n=10 (v0.24.0, FP8 KV, port 18010, 262144 ctx, util=0.49, kv-mem=15360M, seq=64, docker-mem=69g, post-start cgroup enforce)
# Ordering only: text recovery must not start, stop, or restart embedding.
After=network.target vllm-embedding.service
# Direct AEON text endpoint moved behind llm-guard-proxy; wrapper owns stable 18009.
Conflicts=vllm-aeon-27b.service vllm-aeon-27b-dflash-v2.service vllm-aeon-27b-dflash-v3.service vllm-aeon-27b-nvfp4kv.service vllm-aeon-27b-nvfp4kv-test.service vllm-aeon-27b-nvfp4kv-lowmem.service vllm-aeon-27b-nvfp4kv-plus10g.service
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple
Environment=DOCKER_HOST=unix:///run/user/1001/docker.sock
Environment=HF_HUB_OFFLINE=1
Environment=GB10_CGROUP_REGISTRATION_PATH=%t/gb10-memory-guardian/text-cgroup.v1
Environment=GB10_CONTAINER_CIDFILE=%t/gb10-memory-guardian/aeon-text.cid
OOMScoreAdjust=500
# Wrapper/docker-CLI only. Real AEON memory hard-cap is docker --memory 69g plus
# ExecStartPost MemoryMax=69G / MemorySwapMax=0 on docker-*.scope.
# 2026-07-14: 15GiB KV verified 269589 tokens (1.028 full contexts).
# It preserves one full 262144-token context, not two concurrent max requests.
# The first live all-three 3300s window ended with ~31.6GiB MemAvailable,
# zero threshold events, and zero cgroup OOM kills.
# Docker/systemd caps are safety limits, not physical NVML/UMA ceilings.
MemoryMax=256M
MemorySwapMax=0
# Recover from real vLLM/container crashes such as EngineDeadError exits.
# Manual `systemctl --user stop` remains a clean stop and is not restarted.
Restart=on-failure
RestartSec=30
TimeoutStartSec=1200
TimeoutStopSec=60

ExecStartPre=/bin/bash -c 'for i in $(seq 1 120); do curl -fsS --max-time 3 http://100.105.4.92:18012/v1/models >/dev/null 2>&1 && exit 0; sleep 5; done; echo "embedding not ready after 10min"; exit 1'
ExecStartPre=/usr/bin/install -d -m 0700 %t/gb10-memory-guardian
ExecStartPre=/usr/bin/rm -f %t/gb10-memory-guardian/text-cgroup.v1 %t/gb10-memory-guardian/aeon-text.cid
ExecStartPre=-/usr/bin/docker rm -f vllm-aeon-27b-dflash-n12 vllm-aeon-27b-dflash vllm-aeon-27b
ExecStart=/usr/bin/docker run --rm --cidfile=%t/gb10-memory-guardian/aeon-text.cid --name vllm-aeon-27b-dflash-n12 \
  --cgroup-parent app.slice \
  --gpus all --ipc host \
  -p 100.105.4.92:18010:8000 \
  -v /home/obj/.cache/huggingface/hub/models--AEON-7--Qwen3.6-27B-AEON-Ultimate-Uncensored-Multimodal-NVFP4-MTP-XS:/models/aeon-cache:ro \
  -v /home/obj/.cache/huggingface/hub/models--z-lab--Qwen3.6-27B-DFlash:/models/dflash-cache:ro \
  -v /home/obj/scripts/aeon_hang_guard.py:/opt/hang_guard/aeon_hang_guard.py:ro \
  -v /home/obj/scripts/aeon_vllm_wrapper.py:/opt/hang_guard/aeon_vllm_wrapper.py:ro \
  -e HF_HUB_OFFLINE=1 \
  -e PYTHONPATH=/opt/hang_guard \
  -e AEON_DEFAULT_THINKING_TOKEN_BUDGET=32768 \
  -e AEON_THINKING_BUDGET_GUARD=0 \
  --memory 69g \
  --memory-swap 69g \
  --memory-swappiness 0 \
  --oom-score-adj 500 \
  --entrypoint python3 \
  ghcr.io/aeon-7/aeon-vllm-ultimate@sha256:f6d453d0b4a7ef90eefee486f4ff769cc2e1bb1e206df16d70370da09c02203c \
  /opt/hang_guard/aeon_vllm_wrapper.py serve \
    /models/aeon-cache/snapshots/6394b93fa092dcbfbf09e952f56f30337509d17c \
    --host 0.0.0.0 --port 8000 \
    --served-model-name "aeon-ultimate" "qwen3.6-27b-decensor-by-aeon" \
    --quantization modelopt \
    --kv-cache-dtype fp8_e4m3 \
    --attention-backend TRITON_ATTN \
    --mamba-cache-dtype float32 \
    --max-model-len 262144 \
    --max-num-seqs 64 \
    --max-num-batched-tokens 32768 \
    --enable-prefix-caching \
    --gpu-memory-utilization 0.49 \
    --kv-cache-memory-bytes 15360M \
    --speculative-config '{"method":"dflash","model":"/models/dflash-cache/snapshots/0919688658996800f86b895034249700e9481106","num_speculative_tokens":10,"attention_backend":"TRITON_ATTN"}' \
    --enable-chunked-prefill \
    --no-scheduler-reserve-full-isl \
    --reasoning-parser qwen3 \
    --enable-auto-tool-choice \
    --tool-call-parser qwen3_coder \
    --default-chat-template-kwargs '{"enable_thinking": false}' \
    --override-generation-config '{"temperature": 0.6, "top_p": 0.95, "top_k": 20}' \
    --trust-remote-code
# Rootless Docker can accept --memory-swap while leaving memory.swap.max=max.
# Enforce sets docker-*.scope MemoryMax=69G + MemorySwapMax=0 and verifies both.
# Cap arithmetic 69+20+18=107GiB is policy, not guaranteed physical headroom.
ExecStartPost=/home/obj/.local/bin/gb10_enforce_docker_cgroup_limits.sh vllm-aeon-27b-dflash-n12 69
ExecStop=/usr/bin/timeout --signal=TERM --kill-after=5 30 /usr/bin/docker stop --time 20 vllm-aeon-27b-dflash-n12
ExecStopPost=/bin/bash -c 'cidfile=%t/gb10-memory-guardian/aeon-text.cid; [[ -f "$cidfile" && ! -L "$cidfile" ]] || exit 0; IFS= read -r cid <"$cidfile"; [[ "$cid" =~ ^[0-9a-f]{64}$ ]] || exit 1; /usr/bin/timeout --signal=TERM --kill-after=2 3 /usr/bin/docker inspect "$cid" >/dev/null 2>&1 || exit 0; /usr/bin/timeout --signal=TERM --kill-after=2 8 /usr/bin/docker stop --time 5 "$cid" >/dev/null 2>&1 || /usr/bin/timeout --signal=TERM --kill-after=2 3 /usr/bin/docker kill "$cid" >/dev/null 2>&1; running=$(/usr/bin/timeout --signal=TERM --kill-after=2 3 /usr/bin/docker inspect -f "{{.State.Running}}" "$cid" 2>/dev/null || true); [[ -z "$running" || "$running" == false ]]'
ExecStopPost=/usr/bin/rm -f %t/gb10-memory-guardian/text-cgroup.v1 %t/gb10-memory-guardian/aeon-text.cid
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
