# systemd/vllm-aeon-27b-dflash.service
# ── Baseline (clean start 2026-07-16, emb+rr+text co-resident) ──
# GPU: 39.52G  RSS: 6.50G  Startup: 1034s (~17.2min, FlashInfer JIT serial)
# Config: util=0.355 AUTO KV, batch=4096, MAX_JOBS=1, FP8 KV
# KV capacity: 286962 tokens (24818 / 9.47% above one 262144-token context)
# MemAvail after all 3: 39.75G  (cold-ready GPU total 67.89G of 121.6G)
# ─────────────────────────────────────────────────────────────────
[Unit]
Description=vLLM AEON-7 Qwen3.6-27B DFlash n=10 (v0.25.0, seqs=16, batch=4096, util=0.355, AUTO KV, FULL_DECODE_ONLY, FP8 KV, 262144 ctx, multimodal, port 18010)
# Independent of embedding/reranker so text can restart alone (no HTTP poll).
After=network.target
Conflicts=vllm-aeon-27b.service vllm-aeon-27b-dflash-v2.service vllm-aeon-27b-dflash-v3.service vllm-aeon-27b-nvfp4kv.service vllm-aeon-27b-nvfp4kv-test.service vllm-aeon-27b-nvfp4kv-lowmem.service vllm-aeon-27b-nvfp4kv-plus10g.service
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple
Environment=DOCKER_HOST=unix:///run/user/1001/docker.sock
Environment=HF_HUB_OFFLINE=1
Environment=GB10_CGROUP_REGISTRATION_PATH=%t/gb10-memory-guardian/text-cgroup.v1
Environment=GB10_CONTAINER_CIDFILE=%t/gb10-memory-guardian/aeon-text.cid
Environment=AEON_URL=http://100.105.4.92:18010
# Prefer killing text before reranker under host OOM (reranker uses 500).
OOMScoreAdjust=800
# Wrapper only. The container's 128G cap is not a physical NVML/UMA ceiling.
# Memory budget (3-service ~80G+ class): weights~22G + pinned KV 15G + graphs~0.5G
# + emb~20G + rr~18G ≈ 75G engine-side policy, not a cgroup guarantee.
# 2026-07-14: 15GiB KV verified 269589 tokens (1.028 full contexts).
# It preserves one full 262144-token context, not two concurrent max requests.
# The first live all-three 3300s window ended with ~31.6GiB MemAvailable,
# zero threshold events, and zero cgroup OOM kills.
# Auto-KV at util=0.6 grew to ~46G and left only ~10G MemAvailable — too fat.
# AEON-recommended DFlash shape: seqs=16, batch=8192, util=0.6, FULL_DECODE_ONLY.
# NO --kv-cache-memory-bytes: AEON confirmed pinning KV bypasses UMA guard
# (#46932). determine_available_memory early-returns on explicit bytes,
# skipping memory_profiling() + profile_cudagraph_memory() — the only paths
# patched by the UMA free-memory + cudagraph headroom check. AUTO sizing
# lets these clamps bound the pool based on actual free memory.
# batch=4096: reduces profile-run compile/cudagraph peak shape to fit co-resident
# emb+rr without stopping them. Steady-state decode unaffected (governed by seqs);
# only cold long-prompt prefill TTFT increases via more chunked-prefill iterations.
# Prefix-cache/KV-cache hit rate is unaffected (lives inside the pool, not batch-sized).
# util=0.6 is AEON's DFlash-safe GB10 target (0.88 hard ceiling).
# Restart via gb10_restart_text_safe.sh (stops rr first, never emb).
MemoryMax=256M
MemorySwapMax=0
# Recover from any stop cause (crash, OOM kill, guardian cgroup.kill).
# Manual `systemctl --user stop` is NOT restarted because systemd treats
# an explicit stop as a clean deactivation, not a failure.
Restart=always
RestartSec=30
TimeoutStartSec=2400
TimeoutStopSec=60

ExecStartPre=/usr/bin/install -d -m 0700 %t/gb10-memory-guardian
ExecStartPre=/usr/bin/install -d -m 0700 /home/obj/.cache/vllm-compile/aeon-qwen36-v025-18c09e
ExecStartPre=-/usr/bin/docker rm -f vllm-aeon-27b-dflash-n12 vllm-aeon-27b-dflash vllm-aeon-27b
ExecStartPre=/usr/bin/rm -f %t/gb10-memory-guardian/text-cgroup.v1 %t/gb10-memory-guardian/aeon-text.cid
ExecStart=/usr/bin/docker run --rm --cidfile=%t/gb10-memory-guardian/aeon-text.cid --name vllm-aeon-27b-dflash-n12 \
  --cgroup-parent app.slice \
  --gpus all --ipc host \
  -p 100.105.4.92:18010:8000 \
  -v /home/obj/.cache/huggingface/hub/models--AEON-7--Qwen3.6-27B-AEON-Ultimate-Uncensored-Multimodal-NVFP4-MTP-XS:/models/aeon-cache:ro \
  -v /home/obj/.cache/huggingface/hub/models--z-lab--Qwen3.6-27B-DFlash:/models/dflash-cache:ro \
  -v /home/obj/.cache/vllm-compile/aeon-qwen36-v025-18c09e:/var/cache/vllm/aeon-qwen36-v025 \
  -v /home/obj/scripts/aeon_hang_guard.py:/opt/hang_guard/aeon_hang_guard.py:ro \
  -v /home/obj/scripts/aeon_vllm_wrapper.py:/opt/hang_guard/aeon_vllm_wrapper.py:ro \
  -e HF_HUB_OFFLINE=1 \
  -e MAX_JOBS=1 \
  -e CMAKE_BUILD_PARALLEL_LEVEL=1 \
  -e PYTHONPATH=/opt/hang_guard \
  -e AEON_DEFAULT_THINKING_TOKEN_BUDGET=32768 \
  -e AEON_THINKING_BUDGET_GUARD=0 \
  --memory-swappiness 0 \
  --memory 128g \
  --memory-swap 128g \
  --oom-score-adj 800 \
  --entrypoint python3 \
  ghcr.io/aeon-7/aeon-vllm-ultimate@sha256:18c09e6b80141a530285160781f7fa720a78ef91143b3c15a65a8c9641b44e55 \
  /opt/hang_guard/aeon_vllm_wrapper.py serve \
    /models/aeon-cache/snapshots/6394b93fa092dcbfbf09e952f56f30337509d17c \
    --host 0.0.0.0 --port 8000 \
    --served-model-name "aeon-ultimate" "qwen3.6-27b-decensor-by-aeon" \
    --quantization modelopt \
    --kv-cache-dtype fp8_e4m3 \
    --attention-backend TRITON_ATTN \
    --mamba-cache-dtype float32 \
    --max-model-len 262144 \
    --max-num-seqs 16 \
    --max-num-batched-tokens 4096 \
    --max-num-partial-prefills 1 \
    --max-long-partial-prefills 1 \
    --enable-prefix-caching \
    --gpu-memory-utilization 0.355 \
    --enable-chunked-prefill \
    --scheduler-reserve-full-isl \
    --compilation-config "{\"mode\":3,\"cudagraph_mode\":\"FULL_DECODE_ONLY\",\"cudagraph_num_of_warmups\":0,\"cache_dir\":\"/var/cache/vllm/aeon-qwen36-v025\"}" \
    --jit-monitor-mode warn \
    --max-logprobs 20 \
    --speculative-config "{\"method\":\"dflash\",\"model\":\"/models/dflash-cache/snapshots/0919688658996800f86b895034249700e9481106\",\"num_speculative_tokens\":10,\"attention_backend\":\"TRITON_ATTN\"}" \
    --reasoning-parser qwen3 \
    --enable-auto-tool-choice \
    --tool-call-parser qwen3_coder \
    --default-chat-template-kwargs "{\"enable_thinking\":false}" \
    --override-generation-config "{\"temperature\":0.6,\"top_p\":0.95,\"top_k\":20}" \
    --trust-remote-code
# Multimodal retained: no --language-model-only
# Publish the exact Docker cgroup for llm-guard-proxy's integrated guardian.
ExecStartPost=/home/obj/.local/bin/llm_guard_proxy_publish_cgroup_registration.sh
# Raw :18010 readiness (defaults target Guard :18009).
ExecStartPost=/home/obj/.local/bin/gb10_service_ready.sh chat http://100.105.4.92:18010 aeon-ultimate --deadline 2200
ExecStop=/usr/bin/timeout --signal=TERM --kill-after=5 30 /usr/bin/docker stop --time 20 vllm-aeon-27b-dflash-n12
ExecStopPost=/usr/bin/rm -f %t/gb10-memory-guardian/text-cgroup.v1 %t/gb10-memory-guardian/aeon-text.cid
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
