# systemd/vllm-qwen3-reranker-8b.service
[Unit]
Description=vLLM Qwen3-Reranker-8B (AEON v0.25 GB10 image, BF16, 40k ctx, pooling reranker, util=0.22, raw backend port 18013, behind llm-guard-proxy listeners 18009/18003/18005, 64 seq, post-start swap enforcement)
# lifecycle-independent from text: guardian recovery must not change reranker state.
After=network.target
Conflicts=querit-4b-reranker.service vllm-qwen35-hauhau.service
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple
Environment=DOCKER_HOST=unix:///run/user/1001/docker.sock
Environment=HF_HUB_OFFLINE=1
OOMScoreAdjust=800
Restart=on-failure
RestartSec=30
TimeoutStartSec=1320
TimeoutStopSec=60

ExecStartPre=-/usr/bin/docker rm -f vllm-qwen3-reranker-8b
ExecStart=/usr/bin/docker run --rm --name vllm-qwen3-reranker-8b \
  --gpus all --ipc host --dns 8.8.8.8 --dns 1.1.1.1 \
  -p 100.105.4.92:18013:8000 \
  -v /home/obj/.cache/huggingface:/root/.cache/huggingface \
  -e HF_HUB_OFFLINE=1 \
  --memory 24g \
  --memory-swap 24g \
  --memory-swappiness 0 \
  --oom-score-adj 800 \
  --entrypoint python3 \
  ghcr.io/aeon-7/aeon-vllm-ultimate:2026-07-14-v0.25.0 \
  /usr/local/bin/vllm serve Qwen/Qwen3-Reranker-8B \
    --host 0.0.0.0 --port 8000 \
    --served-model-name qwen3-reranker-8b Qwen/Qwen3-Reranker-8B \
    --runner pooling \
    --dtype bfloat16 \
    --max-model-len 40960 \
    --max-num-batched-tokens 40960 \
    --max-num-seqs 64 \
    --gpu-memory-utilization 0.22 \
    --kv-cache-memory-bytes 5820M \
    --enforce-eager \
    --chat-template /vllm-workspace/examples/pooling/score/template/qwen3_reranker.jinja \
    --hf-overrides '{"architectures":["Qwen3ForSequenceClassification"],"classifier_from_token":["no","yes"],"is_original_qwen3_reranker":true}' \
    --trust-remote-code
# Rootless Docker can accept --memory-swap while leaving memory.swap.max=max.
# Enforce the generated docker scope directly. This disabled 24G fallback is not
# part of the ordinary 69+20+18=107GiB AEON+embedding+Querit source profile.
ExecStartPost=/home/obj/.local/bin/gb10_enforce_docker_cgroup_limits.sh vllm-qwen3-reranker-8b 24
ExecStop=/usr/bin/docker stop vllm-qwen3-reranker-8b
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
