# systemd/vllm-embedding.service
[Unit]
Description=vLLM Qwen3-Embedding-8B (AEON v0.24 GB10 image, raw backend port 18012, behind llm-guard-proxy listeners 18009/18002/18005, reliability-first, 40k native max, 64 seq, post-start swap enforcement)
After=network.target
StartLimitIntervalSec=600
StartLimitBurst=3

[Service]
Type=simple
Environment=DOCKER_HOST=unix:///run/user/1001/docker.sock
Environment=HF_HUB_OFFLINE=1
# Embedding is the reliability-critical service. Ideal would be oom_score_adj=0
# but rootless docker (rootlesskit pid 3907540 adj=200) caps the child floor
# at 200 — obj user cannot lower a process's oom_score_adj without sudo. The
# --oom-score-adj 0 flag below is kept as documentation of intent; docker
# silently accepts it but the value remains 200 in HostConfig. Protection
# relies on temporary workloads being launched with --oom-score-adj 1000
# (giving an 800-point gap), so the kernel kills them first under pressure.
# To fully reach adj=0, run `sudo choom -p <rootlesskit_pid> -n 0` then
# `systemctl --user restart vllm-embedding`.
Restart=on-failure
RestartSec=30
TimeoutStartSec=600
TimeoutStopSec=60

ExecStartPre=-/usr/bin/docker rm -f vllm-embedding
ExecStart=/usr/bin/docker run --rm --name vllm-embedding \
  --gpus all --ipc host --dns 8.8.8.8 --dns 1.1.1.1 \
  -p 100.105.4.92:18012:8000 \
  -v /home/obj/.cache/huggingface:/root/.cache/huggingface \
  -e HF_HUB_OFFLINE=1 \
  --memory 24g \
  --memory-swap 24g \
  --memory-swappiness 0 \
  --oom-score-adj 0 \
  --entrypoint python3 \
  ghcr.io/aeon-7/aeon-vllm-ultimate:2026-07-01-v0.24.0 \
  /usr/local/bin/vllm serve Qwen/Qwen3-Embedding-8B \
    --host 0.0.0.0 --port 8000 \
    --served-model-name qwen3-embedding-8b Qwen/Qwen3-Embedding-8B \
    --convert embed \
    --max-model-len 40960 \
    --max-num-batched-tokens 8192 \
    --max-num-seqs 64 \
    --kv-cache-memory-bytes 5820M \
    --enforce-eager
# Rootless Docker can accept --memory-swap while leaving memory.swap.max=max.
# Enforce the generated docker scope directly and verify the 64+24+24=112GiB
# vLLM budget leaves ~9.6GiB headroom on the 121.6GiB GB10 host.
ExecStartPost=/home/obj/.local/bin/gb10_enforce_docker_cgroup_limits.sh vllm-embedding 24
ExecStop=/usr/bin/docker stop vllm-embedding
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
